# .config/autostart/conky.desktop
[Desktop Entry]
Type=Application
Name=conky
Exec=conky --pause=10 -c "/home/emilis/.config/conky/conky_light"
StartupNotify=false
Terminal=false
